# src/spectrum/dashboard/spectrum-os.service
[Unit]
Description=Spectrum OS (光谱 OS)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/SpectrumOS
EnvironmentFile=/opt/SpectrumOS/.env
ExecStart=/usr/bin/python3 -m spectrum.main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
